# source/linux/arch_linux/archlinux_config_ip/dhclient.service
[Unit]
Description= DHCP Client
Before=network.target

[Service]
Type=simple
ExecStart=/sbin/dhclient wlan0 -v
ExecStop=/sbin/dhclient wlan0 -r

[Install]
WantedBy=multi-user.target
