# source/kvm/libvirt/network/libvirt_bridged_network/udev_99-bridge.rules
ACTION=="add", SUBSYSTEM=="module", KERNEL=="bridge", RUN+="/sbin/sysctl -p /etc/sysctl.d/bridge.conf"
